# crontab
* 1 * * * /usr/bin/docker restart kbbi-scraper
